# dots_.config/systemd/user/smartctl-logger.service
[Unit]
Description=Smartctl hourly logger

[Service]
Type=oneshot
ExecStart=/home/kamisaki/.local/bin/smartctl_logger.sh
