# dots_.config/systemd/user/smartctl-logger.timer
[Unit]
Description=Run smartctl logger every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
